# infra/com.danaos.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd plist for the Dana OS daily backup.
  Install once:
    sed -e "s|__REPO_ROOT__|$(pwd)|g" -e "s|__HOME__|$HOME|g" \
      infra/com.danaos.backup.plist \
      > ~/Library/LaunchAgents/com.danaos.backup.plist
    launchctl load ~/Library/LaunchAgents/com.danaos.backup.plist

  Run immediately (for testing):
    launchctl start com.danaos.backup

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.danaos.backup.plist
    rm ~/Library/LaunchAgents/com.danaos.backup.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.danaos.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/backup.sh</string>
  </array>

  <!-- Daily at 02:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/dana-os-backup.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/dana-os-backup.log</string>

  <!-- Run the missed job the next time the machine is awake if 02:00 was skipped. -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
